# systemd/.config/systemd/user/owncloud.service
[Unit]
Description=Owncloud

[Service]
Type=simple
ExecStart=/usr/bin/owncloud
Restart=always

[Install]
WantedBy=default.target
